# platform/broadcom/sonic-platform-modules-accton/wedge100s-32x/service/wedge100s-bmc-poller.timer
[Unit]
Description=Run Wedge 100S-32X BMC sensor daemon every 10 seconds
After=wedge100s-platform-init.service
Requires=wedge100s-platform-init.service

[Timer]
# Fire immediately when this timer unit activates (after platform-init).
OnActiveSec=0
# Subsequent runs every 10 s.  thermalctld polls every 60 s; 10 s
# ensures the daemon output is always <10 s stale.
OnUnitActiveSec=10
AccuracySec=1

[Install]
WantedBy=multi-user.target
