# platform/broadcom/sonic-platform-modules-accton/wedge100s-32x/service/wedge100s-flex-counter-daemon.service
[Unit]
Description=Wedge100S flex-counter daemon (breakout sub-port counters via bcmcmd)
After=syncd.service database.service
Wants=database.service

[Service]
Type=simple
ExecStart=/usr/bin/wedge100s-flex-counter-daemon
Restart=always
RestartSec=5s
TimeoutStopSec=5s

[Install]
WantedBy=multi-user.target
