# platform/broadcom/sonic-platform-modules-accton/wedge100s-32x/service/wedge100s-i2c-poller.service
[Unit]
Description=Wedge100S QSFP presence + EEPROM cache (event-driven)
Documentation=file:///usr/bin/wedge100s-i2c-daemon.c
After=wedge100s-platform-init.service
Requires=wedge100s-platform-init.service

[Service]
Type=oneshot
ExecStart=/usr/bin/wedge100s-i2c-daemon poll-presence
TimeoutStartSec=15
# Suppress routine Starting/Finished lifecycle messages (fires every 3 s).
LogLevelMax=notice
