# platform/broadcom/sonic-platform-modules-accton/wedge100s-32x/etc/systemd/system/swss.service.d/restart-limit.conf
# Wedge100S: Tomahawk ASIC init takes 30-60s. orchagent (inside swss
# container) sends INIT_VIEW to syncd immediately on startup; if syncd
# hasn't finished ASIC init, orchagent aborts with SIGABRT.  The default
# StartLimitBurst=3 with RestartSec=30 exhausts all retries before the
# ASIC is ready.  Raise to 10 attempts over 600s so orchagent survives
# the race without manual intervention.
[Unit]
StartLimitBurst=10
StartLimitIntervalSec=600
